# templates/funkwhale.target.j2
[Unit]
Description=Funkwhale
Wants={{ funkwhale_systemd_service_name }}-server.service {{ funkwhale_systemd_service_name }}-worker.service {{ funkwhale_systemd_service_name }}-beat.service
